# devel/etc/catkin/profile.d/25.turtlebot-stdr.sh
/frontier_expo/frontier_exploration/devel/.private/turtlebot_stdr/etc/catkin/profile.d/25.turtlebot-stdr.sh